# sync.cron
*/30 * * * * /usr/local/bin/pihole-cloudsync --push
